# etman-server/etman-server.service
[Unit]
Description=ET:Legacy ETMan Server (Voice + Sounds + Admin)
After=network.target etserver.service postgresql.service
Wants=etserver.service

[Service]
Type=simple
User=andy
WorkingDirectory=/home/andy/etlegacy
Environment="DATABASE_URL=postgresql://etpanel:etpanel123@localhost:5432/etpanel"
ExecStart=/home/andy/etlegacy/etman_server 27961 27960
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65535
MemoryMax=256M

[Install]
WantedBy=multi-user.target
